# dist/torque-fast-sched.service
[Unit]
Description=MetaCentrum Torque Scheduler
#Requires=local-fs.target named.target network.target pbs_server.target
After=local-fs.target named.target network.target pbs_server.service

[Service]
Type=forking
ExecStart=/usr/sbin/torque-fast-sched
PIDFile=/var/spool/torque/sched_priv/sched.lock
TimeoutStopSec=2400

[Install]
WantedBy=multi-user.target
